# services/frontend-dev/run
#!/usr/bin/env -S bash -e
cd ../../frontend
exec pnpm dev
