# services/backend-dev/run
#!/usr/bin/env -S bash -e
cd ../../backend
exec air
